# check_postgres_replication.cron

0 * * * * root /usr/bin/check_postgres_replication >/dev/null 2>&1
